# bash/.config/bash/profile.d/10-completion.bash
# Use Apple Command Line Tools Git completion directly.
if [ -r "/Library/Developer/CommandLineTools/usr/share/git-core/git-completion.bash" ]; then
  . "/Library/Developer/CommandLineTools/usr/share/git-core/git-completion.bash"
fi
